# root/etc/s6-overlay/s6-rc.d/digrouz-base-adduser/run
#!/command/with-contenv bash

NUID=${UID:-12345}
NGID=${GID:-12345}

groupmod -o -g "${NGID}" appuser
usermod -o -u "${NUID}" appuser
